# src/scrapping/scraping.cron
* */1 * * * python3 scraping_coronavirus.py && npm run build --path /var/OCP/repos/coronaMap