# src/scrapping/scraping.cron
MAILTO=elmaghraoui.ahmed@gmail.com
*/5 * * * * . python3 scraping_coronavirus.py &>buildCron.log && npm run build --path /var/OCP/repos/coronaMap 
