# src/scrapping/scraping.cron
*/5 * * * * . python3 scraping_coronavirus.py >> scraping.log && npm run build --path /var/OCP/repos/coronaMap >> buildCron.log

